# misc/crontab
0 0 * * * /home/mtakagi/.local/bin/fix-default-browser
